# root/etc/services.d/rhodecode/run
#!/usr/bin/with-contenv bash
exec /root/.rccontrol-profile/bin/supervisord -c /root/.rccontrol/supervisor/supervisord.ini
